# gunicorn.socket
[Unit]
Description=gunicorn socket

[Socket]
ListenStream=/path_to/gunicorn.sock

[Install]
WantedBy=sockets.target
